# tests/fixtures/scheduler/mybench-scan.service
# Managed by mybench; do not edit.
[Unit]
Description=Run the private mybench reconciliation scan

[Service]
Type=oneshot
ExecStart="/opt/mybench/bin/mybench" scan --quiet --scheduled
NoNewPrivileges=true
PrivateTmp=true
